# deploy/cloudflare-tunnel/com.acceptme.amt-bio-digest-tunnel.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.acceptme.amt-bio-digest-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/ABSOLUTE/PATH/TO/bio-literature-digest-web/bio-literature-config/tunnel/web/config.yml</string>
    <string>run</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/bio-literature-digest-web</string>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/bio-literature-digest-web/bio-literature-config/runtime/web/amt-bio-digest-tunnel.log</string>

  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/bio-literature-digest-web/bio-literature-config/runtime/web/amt-bio-digest-tunnel.log</string>
</dict>
</plist>
